# scripts/samo-server.service
[Unit]
Description=Samo Server — self-hosted music, audiobooks, podcasts, radio
Documentation=https://github.com/bouliehaan/samo-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=samo
Group=samo
WorkingDirectory=/var/lib/samo
ExecStart=/opt/samo/samo-server
Restart=on-failure
RestartSec=5s

# Defaults. Override these by dropping a file in /etc/systemd/system/samo-server.service.d/
Environment=SAMO_DATA_DIR=/var/lib/samo
Environment=SAMO_ADDR=:6969
Environment=SAMO_SCAN_ON_START=false

# Light hardening. The operator can tighten this with a drop-in if their
# media folders aren't under common locations like /srv, /mnt, /media.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/samo
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Journal is the easiest path for log discovery on a stock Ubuntu box.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
